# shotover.service
[Unit]
Description=ShotoverProxy

Wants=network.target
After=syslog.target network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/shotover-proxy -t /etc/shotover/config/topology.yaml -c /etc/shotover/config/config.yaml
Restart=on-failure
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target

