# shotover.service
[Unit]
Description=ShotoverProxy

Wants=network.target
After=syslog.target network-online.target

[Service]
Type=simple
ExecStart=/root/shotover/shotover-proxy
Restart=on-failure
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target

